# supervisord/nordvpn.conf
[program:nordvpn]
command=/nordvpn_setup.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/nordvpn.err.log
stdout_logfile=/var/log/nordvpn.out.log
priority=10
